# supervisord.conf
[program:ircbot]
command = ruby /home/dotcloud/current/lib/shewbot.rb

